# supervisor/cunxun.ini
[program:cunxun]
command=/data/apps/go/wangqingang/cunxun/bin/cunxun --config conf/config.toml
directory=/data/apps/go/wangqingang/cunxun/
autostart=true
autorestart=true
startsecs=2
startretries=3
stdout_logfile=/data/logs/supervisor/cunxun/access.log
stdout_logfile_maxbytes=100MB
stdout_logfile_backups=20
stderr_logfile=/data/logs/supervisor/cunxun/stderr.log
stderr_logfile_maxbytes=100MB
stderr_logfile_backups=2
